# COMe/come_heartbeat_server@.service
[Unit]
Description=COMe Heartbeat Server

[Service]
ExecStart=/opt/fungible/etc/come_heartbeat_server.py
StandardInput=socket
